# 90-vpdupdate.rules
KERNELS=="*", ACTION=="*", RUN+="/bin/touch /var/lib/lsvpd/run.vpdupdate"
